# docker/aio/program-redis.conf
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis-aio.conf
user=redis
priority=20
autostart=true
autorestart=true
startsecs=3
# SIGTERM triggers a graceful shutdown that flushes the AOF before exiting.
stopsignal=TERM
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
