# docker/aio/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ----------------------------------------------------------------------------
# Shutdown ordering: lower priority -> started first, stopped LAST.
#   docmost (30) -> stopped first  (drains HTTP, closes DB pool)
#   redis   (20) -> stopped second
#   postgres(10) -> stopped last   (final checkpoint after no clients remain)
# ----------------------------------------------------------------------------

[program:postgres]
command=/usr/lib/postgresql/%(ENV_PG_VERSION)s/bin/postgres -D %(ENV_PG_DATA)s -c config_file=%(ENV_PG_DATA)s/postgresql.conf
user=postgres
priority=10
autostart=true
autorestart=true
startsecs=5
# SIGINT = "fast shutdown": cancels active queries, performs a checkpoint,
# flushes WAL, then exits. SIGTERM (smart shutdown) waits indefinitely for
# clients to disconnect, which would block container shutdown.
stopsignal=INT
stopwaitsecs=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redis]
command=/usr/bin/redis-server /etc/redis/redis-aio.conf
user=redis
priority=20
autostart=true
autorestart=true
startsecs=3
# SIGTERM triggers a graceful shutdown that flushes the AOF before exiting.
stopsignal=TERM
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:docmost]
command=/usr/local/bin/aio-start-app.sh
user=node
directory=/app
priority=30
autostart=true
autorestart=true
startsecs=10
stopsignal=TERM
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
